# Deployment/TESTING/CRON/NWP_Data/nwp_retrieval.cron
# Defintions environment for cron
SMACK_DIR=/usr/local/smack
SMACK_DIR_BIN=/usr/local/smack/bin
SMACK_DIR_LOG=/usr/local/smack/log
SMACK_DIR_SKEL=/usr/local/smack/skel
SMACK_DIR_TMP=/usr/local/smack/tmp
SMACK_LOAD=/usr/local/smack/log/smack_loaded
SMACK_INSTALL_LOG=/usr/local/smack/log/install_log
CRON_PATH=/usr/local/smack/cron
SHINY_SRV=/srv/shiny-server
API_SRV=/srv/api-server
KEYSTONE_URL="https://keystone-yyc.cloud.cybera.ca:5000/v2.0"
NOVA_URL="https://nova-yyc.cloud.cybera.ca:8774/v2/2b86ecd5b18f4fafb1d55adb79072def"
CINDER_URL="https://cinder-yyc.cloud.cybera.ca:8776/v1/2b86ecd5b18f4fafb1d55adb79072def"
CINDER2_URL="https://cinder-yyc.cloud.cybera.ca:8776/v2/2b86ecd5b18f4fafb1d55adb79072def"
GLANCE_URL="http://glance-yyc.cloud.cybera.ca:9292"
EC2_URL="https://nova-yyc.cloud.cybera.ca:8773/services/Cloud"
SWIFT_URL="https://swift-yyc.cloud.cybera.ca:8080/v1/AUTH_2b86ecd5b18f4fafb1d55adb79072def"
OS_PROJECT_NAME="SMACK"
OS_ZONE="Nova"
OS_REGION="Calgary"
STORAGE_ACCT="AUTH_4b6be558d44e4dba8fb6e4aa49934c0b"
STORAGE_TOKEN="7eefd48208754002a2e03bf0de11c3e4"
STORAGE_URL="https://swift-yyc.cloud.cybera.ca:8080/v1/AUTH_4b6be558d44e4dba8fb6e4aa49934c0b"
PATH=/bin:/usr/bin:/usr/local/smack/bin:/usr/local/smack/cron/bin:/usr/local/bin

# Retrieve NWP Data (Every 6 Hours)
5 */6 * * * /usr/local/smack/cron/bin/ret_nwp.sh
# Check NWP Data (Every 6 Hours)
25 */6 * * * /usr/local/smack/cron/bin/chk_nwp.sh
# Store NWP Data (Every 6 Hours)
45 */6 * * * /usr/local/smack/cron/bin/str_nwp.sh
# Clear NWP Data (Daily)
0 0 * * * /usr/local/smack/cron/bin/clr_nwp.sh

## HISTORICAL POWER DATA
* * * * * /usr/local/smack/cron/bin/hist.sh